# docker/supervisor/dnsmasq.ini
[program:dnsmasq]
command=dnsmasq -d
user=root

[program:watchdog]
command=/var/www/html/docker/watchdog.sh
user=root
